# timers/snap.timer
[Unit]
Description=Monitors snapshot allocation

[Timer]
OnCalender=OnCalendar=*:0/15
Persistent=true
