# timers/snap.service
[Unit]
Description=Monitors snapshot allocation

[Service]
Persistent=true
Type=oneshot
ExecStart=/home/jaryd/scratch/arch/timers/monitor.snap.pl
