# timers/backup.home/backup.home.service
[Unit]
Description=backup.home

[Service]
ExecStart=/usr/local/bin/backup.home
