# timers/backup.rotate/backup.rotate.service
[Unit]
Description=backup.rotate

[Service]
ExecStart=/usr/local/bin/tsar -v
