# timers/backup/backup.root.service
[Unit]
Description=backup.root

[Service]
ExecStart=/usr/local/bin/backup.root
